# src/org.tizen.bp.manager.service.in
[Unit]
Description=This is from /usr/share/system.d/system-service/org.tizen.bp.manager.service for DBUS

[D-BUS Service]
Name=org.tizen.bp.manager
Exec=/bin/false
SystemdService=bp-manager.service
